# python/zk-sync/com.shrms.zksync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shrms.zksync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gabe/Herd/Smart-HRMS/python/zk-sync/.venv/bin/python</string>
        <string>/Users/gabe/Herd/Smart-HRMS/python/zk-sync/sync_agent.py</string>
        <string>--daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gabe/Herd/Smart-HRMS/python/zk-sync</string>

    <key>StandardOutPath</key>
    <string>/Users/gabe/Herd/Smart-HRMS/python/zk-sync/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gabe/Herd/Smart-HRMS/python/zk-sync/sync.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
